# .gomodcache/github.com/awslabs/amazon-eks-ami/nodeadm@v0.0.0-20240229193347-cfab22a10647/test/e2e/infra/systemd/kubelet.service
[Unit]
Description=Run the mock kubelet

[Service]
ExecStart=/usr/bin/kubelet

[Install]
WantedBy=multi-user.target
